# 2017-11-1 单5G/Relay/openwrt35/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

source /etc/profile
/root/WSOL_Relay.sh > /tmp/WSOL_Relay.out 2>&1 &

crontab -u root /root/cron_Relay.config
/etc/init.d/cron start

exit 0
